# roles/wallabag/templates/wallabag-backup.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Run wallabag backup daily

[Timer]
OnCalendar=daily UTC

[Install]
WantedBy=timers.target
