# roles/wallabag/templates/wallabag-backup.service.j2
# {{ ansible_managed }}

[Unit]
Description=Wallabag backup service

[Service]
ExecStart=/home/{{ ansible_user }}/wallabag/wallabag-backup.sh
