# roles/focalboard/templates/focalboard-backup.service.j2
# {{ ansible_managed }}

[Unit]
Description=Focalboard backup service

[Service]
ExecStart=/home/{{ ansible_user }}/focalboard/focalboard-backup.sh
