# roles/focalboard/templates/focalboard-backup.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Run focalboard backup daily

[Timer]
OnCalendar=daily UTC

[Install]
WantedBy=timers.target
